# deploy/systemd/sv-krishna.service
[Unit]
Description=SV-Krishna Offline Voice Assistant
Wants=network-online.target
After=network-online.target

# Optional dependencies (SV-Krishna should still run if these are down).
Wants=ollama.service sv-krishna-whisper.service
After=ollama.service sv-krishna-whisper.service

[Service]
Type=simple
WorkingDirectory=/opt/SV-Krishna
EnvironmentFile=/etc/sv-krishna/sv-krishna.env

# Prefer running the compiled build for responsiveness.
ExecStart=/usr/bin/node /opt/SV-Krishna/dist/index.js

Restart=always
RestartSec=1
TimeoutStopSec=10
KillSignal=SIGINT

# Hardening (relax if you need device access)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

