# deploy/systemd/svkrishna-whisper.service
[Unit]
Description=SV-Krishna Whisper HTTP (faster-whisper)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/svkrishna/app
EnvironmentFile=/opt/svkrishna/app/.env
ExecStart=/opt/svkrishna/venvs/whisper/bin/python /opt/svkrishna/app/python/whisper_http.py

Restart=always
RestartSec=1
TimeoutStopSec=10

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

